# resources/block_trespassers/service
#Destination file: /etc/init.d/update_known_devices
#!/bin/sh /etc/rc.common

readonly USE_PROCD=1
readonly START=99

start_service() {
    source /root/OpenWrtTools/src/logger_helper.sh "update_known_devices"

    local set_name="known_devices"
    local mac_addresses="$( uci show dhcp | grep \.mac= | sed "s/.*mac='\(.*\)'/\1/" | tr "\n" "," | tr " " "," | sed "s/,$//g" )"

    if [ -z "$mac_addresses" ]; then
        log_error "Aborted updating nft set \"$set_name\" with empty list."
        exit 1
    fi

    function update_in_runtime() {
        local family="inet"
        local table="fw4"
        [ $( nft list sets $family | grep -c $set_name) -le 0 ] \
            && nft add set $family $table $set_name { type ether_addr \; }
        nft flush set $family $table $set_name
        nft add element $family $table $set_name { $mac_addresses }
    }; update_in_runtime

    function update_static_configuration() {
        local new_line="elements = { $mac_addresses }"
        sed -i "s/#\?elements.*=.*{.*}/$new_line/" "$SET_FILE"
    }; update_static_configuration

    log_info "Updated nft set \"$set_name\"."
}

service_triggers() {
    procd_add_reload_trigger "dhcp"
}
